# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for superlists-staging.ottg.eu"
 
start on net-device-up
stop on shutdown
 
respawn
 
setuid borygo
chdir /home/borygo/sites/superlists-staging.borygo.koding.io/source
 
exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/superlists-staging.borygo.koding.io.socket \
    superlists.wsgi:application 
